# debian/snorkel.collector.upstart
description "snorkel_collector"
author      "rdio"

start on startup
stop on shutdown

script
    chdir /opt/snorkel
    exec sudo -u nobody /opt/snorkel/scripts/start_collector
end script

